# deploy/systemd/avelren-postgres-backup.service
[Unit]
Description=Avelren PostgreSQL encrypted backup
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/avelren/backup/backup.env
ExecStart=/usr/local/libexec/avelren-postgres-backup
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
RuntimeDirectory=avelren
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
ReadWritePaths=/var/lib/avelren-backup /var/log
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077
