# deploy/systemd/avelren-postgres-backup.timer
[Unit]
Description=Daily Avelren PostgreSQL backup timer

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=30m
Persistent=true
Unit=avelren-postgres-backup.service

[Install]
WantedBy=timers.target
